# testdata/supervisord.ini
[program:hrafd1]
autostart=false
command = hraftd -rjoin :11001,:11002,:11003,:11004 -haddr :11001
stdout_logfile=/tmp/hraftd1.log
redirect_stderr=true

[program:hrafd2]
autostart=false
command = hraftd -rjoin :11001,:11002,:11003,:11004 -haddr :11002
stdout_logfile=/tmp/hraftd2.log
redirect_stderr=true

[program:hrafd3]
autostart=false
command = hraftd -rjoin :11001,:11002,:11003,:11004 -haddr :11003
stdout_logfile=/tmp/hraftd3.log
redirect_stderr=true

[program:hrafd4]
autostart=false
command = hraftd -rjoin :11001,:11002,:11003,:11004 -haddr :11004
stdout_logfile=/tmp/hraftd4.log
redirect_stderr=true

[inet_http_server]
port=127.0.0.1:9001
username=admin
password=admin

[supervisord]
